# temp/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; 사용자 접근 가능한 소켓 파일 경로

[supervisord]
logfile=/user/viststorage/mlops/supervisor_logs/supervisord.log
pidfile=/user/viststorage/mlops/supervisor_logs/supervisord.pid
childlogdir=/user/viststorage/mlops/supervisor_logs

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[include]
						files = /user/viststorage/mlops/supervisor_conf/programs/*.conf



[program:app1]
command=/user/viststorage/mlops/supervisor_env/bin/python /user/viststorage/mlops/model_list/app1/main.py --host 0.0.0.0 --port 8001
directory=/user/viststorage/mlops/model_list/app1
autostart=false
autorestart=false
stdout_logfile=/user/viststorage/mlops/supervisor_logs/app1.out.log
stderr_logfile=/user/viststorage/mlops/supervisor_logs/app1.err.log

